# dot_config/systemd/user/mako.service.d/override.conf
[Unit]
PartOf=sway-session.target
After=sway-session.target

[Install]
WantedBy=sway-session.target